# hrwros_ws/devel/etc/catkin/profile.d/25.turtlebot-gazebo.sh
/home/venom/ros/hrwros_ws/devel/.private/turtlebot_gazebo/etc/catkin/profile.d/25.turtlebot-gazebo.sh