# hrwros_ws/devel/etc/catkin/profile.d/25.turtlebot-stdr.sh
/home/venom/ros/hrwros_ws/devel/.private/turtlebot_stdr/etc/catkin/profile.d/25.turtlebot-stdr.sh