# hrwros_ws/devel/etc/catkin/profile.d/25.turtlebot-stage.sh
/home/venom/ros/hrwros_ws/devel/.private/turtlebot_stage/etc/catkin/profile.d/25.turtlebot-stage.sh